# crontab
*/2 * * * * /root/automation_project/scrits/auto_backup.sh
